# openwrt-builder/files/etc/crontabs/root
30 0 * * * curl -x http://@secret:CLASH_PROXY_USERNAME:@secret:CLASH_PROXY_PASSWORD@127.0.0.1:10081 -sSL https://github.com/lirundong/homelab-infra/releases/download/sing-box-daemon/config.json.gpg | gpg --batch --passphrase="@secret:MASTER_PASSWORD" --yes -o /root/.config/sing-box/config.json -d && /etc/init.d/sing-box reload
*/15 * * * * PASSWORD="@secret:MASTER_PASSWORD" /usr/bin/timeout -k 3m 1m /usr/bin/python3 /root/util-cookbook/tencent-cloud/register_dns.py -d @secret:DOMAIN -s @secret:DDNS_HOSTNAME -i pppoe-wan --get-ipv4-method taobao --get-ipv6-method ifaddr
