# openwrt-builder/files/etc/uci-defaults/99-install-common.sh
#!/bin/sh

/usr/bin/python3 -m pip install -e /root/common
