# openwrt-builder/files/etc/uci-defaults/99-python-requirements.sh.skip
#!/bin/sh

/usr/bin/python3 -m pip config set global.index-url https://mirrors.tuna.tsinghua.edu.cn/pypi/web/simple
/usr/bin/python3 -m pip install @include:requirements.txt! 
